# playbooks/infrastructure/templates/mnt-nfs.automount.j2
{#
  systemd automount unit for the NFS mount. Used on both the server node
  (self-mount) and client nodes.

  Explicit unit with DefaultDependencies=no and WantedBy=remote-fs.target to
  prevent the fstab-generator from placing it in local-fs.target, which causes
  ordering cycles at boot. See mnt-nfs.mount.j2 for full details.
  (systemd upstream bugs #17657 and #34164 — unfixed as of 2026.)
#}
[Unit]
Description={{ nfs_mount_point }} automount
DefaultDependencies=no
After=network-online.target

[Automount]
Where={{ nfs_mount_point }}
TimeoutIdleSec=0

[Install]
WantedBy=remote-fs.target
