# playbooks/infrastructure/templates/mnt-nfs.mount.j2
{#
  systemd mount unit for the NFS mount. Used on both the server node
  (self-mount) and client nodes.

  On the server node, fstab-based NFS mounts with x-systemd.automount create a
  dependency cycle: nfs-server depends on nfs-mountd, which depends on
  local-fs.target, which pulls in the automount unit (placed there by the
  fstab-generator), which needs nfs-server to succeed. systemd breaks the cycle
  by dropping nfs-mountd, preventing the server from starting.
  (systemd upstream bugs #17657 and #34164 — unfixed as of 2026.)

  On client nodes, fstab-based mounts with x-systemd.automount and an
  After=network-online.target drop-in create a different cycle:
  fstab-generator sets DefaultDependencies=yes, which adds Before=local-fs.target
  to the automount. Since local-fs.target precedes network-online.target in the
  boot sequence, the two ordering constraints are contradictory. systemd resolves
  the deadlock by dropping socket units that share the cycle — including
  iscsid.socket — which prevents Longhorn volumes from attaching.

  Explicit unit files with DefaultDependencies=no resolve both cycles by
  preventing the fstab-generator from pulling the automount into local-fs.target.
  After=nfs-kernel-server.service is included only on the server node, where it
  is safe because the cycle path no longer exists.
#}
[Unit]
Description={{ nfs_mount_point }}
DefaultDependencies=no
After=network-online.target
{% if nfs_server | default(false) %}
After=nfs-kernel-server.service
Wants=nfs-kernel-server.service
{% endif %}

[Mount]
What={{ _nfs_server_host }}:{{ nfs_export_path }}
Where={{ nfs_mount_point }}
Type=nfs
Options=soft,timeo=600,nfsvers=4
TimeoutSec=15

[Install]
WantedBy=remote-fs.target
